# etc/systemd/system/timeshift.timer
# /etc/systemd/system/timeshift.timer
[Unit]
Description=Run 'timeshift --check' at the start of every hour

[Timer]
# Run exactly at the start of every hour
OnCalendar=hourly
# If the machine was turned off during a scheduled slot, run it immediately on boot
# Persistent=true
# Targets the snapshot service explicitly
Unit=timeshift.service

[Install]
WantedBy=timers.target
